# production/supervisor/celerybeat.conf
; ================================
;  celery beat supervisor example
; ================================

[program:celerybeat]
; Set full path to celery program if using virtualenv
command=celery beat -A django_example --loglevel=INFO --pidfile="/tmp/celerybeat.pid" --schedule="/tmp/celerybeat.schedule"

; remove the -A myapp argument if you are not using an app instance

directory=/django_example
user=nobody
numprocs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled = true
stderr_events_enabled = true
autostart=%(ENV_MASTER_WORKER)s
autorestart=%(ENV_MASTER_WORKER)s
startsecs=10