# production/supervisor/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=%(ENV_WEB_SERVER)s
autorestart=%(ENV_WEB_SERVER)s